# platform/nginx-php/docker/config/crontab
* * * * * cd /var/www && php [script - command] >> /dev/stdout 2>&1

